# nanobot.service
[Unit]
Description=Nanobot AI Assistant
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/nanobot
ExecStart=/usr/local/bin/nanobot gateway
Restart=always
RestartSec=10
StandardOutput=append:/root/.nanobot/workspace/logs/service.log
StandardError=append:/root/.nanobot/workspace/logs/service-error.log

# Environment variables (if needed)
# Environment="NANOBOT_CONFIG=/path/to/config.json"

[Install]
WantedBy=multi-user.target
